# files/supervisord/supervisord.conf
[supervisord]
nodaemon=false

[program:postgres]
user=root
group=root
command=/root/postgresql/start.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
redirect_stderr=false
stopsignal=QUIT